# rootfile/rootfs/var/lib/dpkg/info/bluez-alsa-utils.list
/.
/etc
/etc/alsa
/etc/alsa/conf.d
/etc/alsa/conf.d/20-bluealsa.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/bluealsa.conf
/etc/default
/etc/default/bluez-alsa
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bluez-alsa.service
/usr
/usr/bin
/usr/bin/bluealsa
/usr/bin/bluealsa-aplay
/usr/share
/usr/share/doc
/usr/share/doc/bluez-alsa-utils
/usr/share/doc/bluez-alsa-utils/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/bluealsa-aplay.1.gz
/usr/share/man/man8
/usr/share/man/man8/bluealsa.8.gz
/usr/share/doc/bluez-alsa-utils/changelog.Debian.gz
